# pro/HW/get-kernel-logs.timer
[Unit]
Description=Run get-kernel-logs every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
